# deploy/challenge-deploy.service
[Unit]
Description=Auto-deploy plank-challenge from GitHub
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/challenge-deploy.sh
