# deploy/challenge-deploy.timer
[Unit]
Description=Poll GitHub and auto-deploy plank-challenge every 2 min

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
Unit=challenge-deploy.service

[Install]
WantedBy=timers.target
